# scripts/crontab
* * * * * cd /usr/share/nginx/html && php artisan schedule:run